# kbd_idle.service
[Unit]
Description=Keyboard backlight

[Service]
ExecStart=/usr/bin/python /home/hadrien/scripts/kbd_idle.py

[Install]
WantedBy=multi-user.target
